# tools/launchd/com.ae.codex-archaeology-aor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ae.codex-archaeology-aor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tito/lWORKSPACEl/research/neural-memory/tools/codex_always_on_runner.sh</string>
    </array>

    <!-- ALWAYS-ON (replaces daily 04:00 cron) -->
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/tito/.neural_memory/logs/aor-codex-archaeology.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tito/.neural_memory/logs/aor-codex-archaeology.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/Codex.app/Contents/Resources</string>
        <key>HOME</key>
        <string>/Users/tito</string>
        <key>NM_LANE</key>
        <string>cron</string>
        <key>AOR_NAME</key>
        <string>codex-archaeology</string>
        <key>AOR_SCRIPT</key>
        <string>tools/codex_archaeology_pick.sh</string>
        <!-- 60s poll, 30min min, 2hr max — picks random module on fire -->
        <key>AOR_POLL</key>
        <string>60</string>
        <key>AOR_MIN_INTERVAL</key>
        <string>1800</string>
        <key>AOR_MAX_INTERVAL</key>
        <string>7200</string>
    </dict>
</dict>
</plist>
