# tools/launchd/com.ae.neural-reconciliation-reviewer-aor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ae.neural-reconciliation-reviewer-aor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tito/lWORKSPACEl/research/neural-memory/tools/codex_always_on_runner.sh</string>
    </array>

    <!-- ALWAYS-ON (replaces 60min cron) -->
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/tito/.neural_memory/logs/aor-reconciliation-reviewer.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tito/.neural_memory/logs/aor-reconciliation-reviewer.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/Codex.app/Contents/Resources</string>
        <key>HOME</key>
        <string>/Users/tito</string>
        <key>NM_LANE</key>
        <string>nm-builder</string>
        <key>AOR_NAME</key>
        <string>reconciliation-reviewer</string>
        <key>AOR_SCRIPT</key>
        <string>tools/reconciliation_reviewer.sh</string>
        <!-- 60s poll, 30min min between fires, 90min max -->
        <key>AOR_POLL</key>
        <string>60</string>
        <key>AOR_MIN_INTERVAL</key>
        <string>1800</string>
        <key>AOR_MAX_INTERVAL</key>
        <string>5400</string>
    </dict>
</dict>
</plist>
